# deploy/supervisord.conf
; Callendulla — supervisord для запуска api + bot + scheduler в одном
; контейнере. Для production с большой нагрузкой выноси в отдельные
; сервисы через docker-compose.override.yml.

[supervisord]
nodaemon=true
user=callendulla
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; Все три процесса логируют в stdout/stderr — docker logs подхватывает
; их автоматически, никаких отдельных файлов.

[program:api]
command=uvicorn src.api.main:app --host 0.0.0.0 --port 8290 --workers 2 --no-access-log
directory=/app
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:bot]
command=python -m src.bot.main
directory=/app
autostart=true
autorestart=true
startretries=3
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:scheduler]
command=python -m src.scheduler.main
directory=/app
autostart=true
autorestart=true
startretries=3
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30

[group:callendulla]
programs=api,bot,scheduler
